# prometheus-to-logstash/prometheus-to-logstash.service
# /etc/systemd/system/prometheus-to-logstash.service

[Unit]
Description=Prometheus to Logstash Exporter
After=network.target logstash.service
Requires=logstash.service

[Service]
Type=simple
User=monitoring
EnvironmentFile=/etc/default/prometheus-to-logstash
ExecStartPre=/usr/bin/wait-for-it -t 300 ${LOGSTASH_HOST}:${LOGSTASH_PORT}
ExecStartPre=/usr/bin/wait-for-it -t 300 ${PROM_HOST}:${PROM_PORT}
ExecStart=/usr/bin/python3 /usr/local/bin/prometheus-to-logstash.py \
    --prometheus-url http://${PROM_HOST}:${PROM_PORT}/metrics/${TARGETS} \
    --logstash-url ${LOGSTASH_HOST}:${LOGSTASH_PORT} \
    ${ADDITIONAL_OPTIONS}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
